# deploy/telegram-local-agent.service
[Unit]
Description=Telegram Local Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/telegram-local-agent
EnvironmentFile=/etc/telegram-local-agent.env
ExecStart=/usr/local/bin/telegram-local-agent
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/telegram-local-agent/.telegram-memory

[Install]
WantedBy=multi-user.target
