# deploy/telegram-local-agent.service
[Unit]
Description=Telegram Local Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=telegram-local-agent
Group=telegram-local-agent
UMask=0077
WorkingDirectory=/opt/telegram-local-agent
EnvironmentFile=/etc/telegram-local-agent.env
ExecStart=/usr/local/bin/telegram-local-agent
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
ReadWritePaths=/opt/telegram-local-agent/.telegram-memory /opt/telegram-local-agent/.telegram-state

[Install]
WantedBy=multi-user.target
